# bin/systemd/gunicorn-live.socket
[Unit]
Description=gunicorn socket (LIVE)

[Socket]
ListenStream=/run/gunicorn-live.sock

[Install]
WantedBy=sockets.target

